# packaging/ruskery.service
[Unit]
Description=ruskery container registry
Documentation=https://github.com/jacobsamro/ruskery
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=ruskery
Group=ruskery
ExecStart=/usr/local/bin/ruskery --config /etc/ruskery/config.toml serve
Restart=on-failure
RestartSec=2s
# Let the binary bind 80/443 without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/var/lib/ruskery
StateDirectory=ruskery
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
